# playerd-alex/files/playerd.init
#!/bin/sh /etc/rc.common

START=99
STOP=35
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/start_playerd.sh
	procd_set_param respawn 1 1 10000
	procd_close_instance
}

service_triggers()
{
	procd_add_reload_trigger "dhcp"
}

